# pi/systemd/control-server.service
[Unit]
Description=RP3 Tank control server (camera stream + WebSocket + UART)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/opt/rp3tank/pi
ExecStart=/opt/rp3tank/pi/venv/bin/python -m src.main
Restart=always
RestartSec=3

# Allow clean shutdown (STOP to Pico, stop camera)
KillSignal=SIGTERM
TimeoutStopSec=10

[Install]
WantedBy=multi-user.target
